# service/shadowbox.service
[Unit]
Description=Shadowbox RNBO Hardware UI
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/shadowbox
ExecStart=/home/pi/shadowbox/.venv/bin/python -m shadowbox.shadowbox
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/shadowbox

[Install]
WantedBy=multi-user.target
